# ol7-repo-sync.cron
#ol7 repo sync
# @daily root /usr/local/bin/ol7-repo-sync.bash > /dev/null 2>&1
0 3 * * * root /usr/local/bin/ol7-repo-sync.bash > /var/log/ol-repo-sync.log 2>&1
